# supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Quotes around values are not supported, except in the case of
;    the environment= options as shown below.
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".
;  - Command will be truncated if it looks like a config file comment, e.g.
;    "command=bash -c 'foo ; bar'" will truncate to "command=bash -c 'foo ".

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file
chmod=0744                 ; socket file mode (default 0700)

[inet_http_server]         ; inet (TCP) server disabled by default
port=localhost:9001        ; ip_address:port specifier, *:port for all iface

[supervisord]
logfile=/logs/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=10MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/logs/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
childlogdir=/logs            ; 'AUTO' child log dir, default $TEMP
environment=RECORDER_DATA_DIR="/tmp"

; Need this true, or tailing sometimes fails. See https://github.com/Supervisor/supervisor/issues/85.
strip_ansi=true            ; strip ansi escape codes in logs; def. false

; The supervisorctl section configures how supervisorctl will connect to
; supervisord.  configure it match the settings in either the unix_http_server
; or inet_http_server section.
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:sshd]
command=/usr/sbin/sshd -D

[program:asv_bridge]
command=/bin/bash -c "source /opt/mr/ros2_ws/install/setup.bash &&  ros2 run asv_bridge asv_bridge"
priority=0
redirect_stderr=true                    ; redirect proc stderr to stdout (default false)
stdout_logfile=/logs/asv_bridge.log     ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=10MB            ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10

[program:recorder]
command=/bin/bash -c "source /opt/mr/ros2_ws/install/setup.bash && RECORDER_DATA_DIR=/tmp ros2 run recorder recorder"
priority=0
redirect_stderr=true                ; redirect proc stderr to stdout (default false)
stdout_logfile=/logs/recorder.log   ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=10MB        ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10

[program:status]
command=/bin/bash -c "source /opt/mr/ros2_ws/install/setup.bash && ros2 run status status_node"
priority=0
redirect_stderr=true                ; redirect proc stderr to stdout (default false)
stdout_logfile=/logs/status.log   ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=10MB        ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10

[program:pinger]
command=/bin/bash -c "source /opt/mr/ros2_ws/install/setup.bash && ros2 run pinger pinger_node"
priority=0
redirect_stderr=true                ; redirect proc stderr to stdout (default false)
stdout_logfile=/logs/pinger.log   ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=10MB        ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10

[program:mavproxy]
command=mavproxy.py
    --master=/dev/ttyACM0,115200
    --load-module='GPSInput,DepthOutput'
    --source-system=200
    --cmd="set heartbeat 0"
    --out udpin:localhost:9000
    --out udpout:localhost:9002
    --out udpin:0.0.0.0:14660
    --out udpbcast:10.1.1.255:14550
    --mav20
    --aircraft telemetry
    --streamrate 10
priority=0
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/logs/mavproxy.log        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=10MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10